# install/launchd/com.nco.fleet-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.nco.fleet-watchdog</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/nova-ai/.claude/hooks/fleet-session-watchdog.sh</string>
  </array>
  <key>StartInterval</key><integer>600</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/nova-ai/.claude/fleet-watchdog.out</string>
  <key>StandardErrorPath</key><string>/Users/nova-ai/.claude/fleet-watchdog.err</string>
</dict>
</plist>
